# borg/.config/systemd/user/borg-backup.timer
[Unit]
Description=borg backup

[Timer]
Persistent=true
OnCalendar=hourly
Unit=borg-backup.service

[Install]
WantedBy=timers.target
