# borg/.config/systemd/user/borg-prune.timer
[Unit]
Description=prune borg snapshots

[Timer]
Persistent=true
OnCalendar=daily
Unit=borg-prune.service

[Install]
WantedBy=timers.target
