# borg/.config/systemd/user/borg-backup.service
[Unit]
Description=borg backup
Before=borg-prune.service

[Service]
Type=oneshot
ExecStart=%h/.bin/borg_backup

[Install]
WantedBy=multi-user.target
