# restic/.config/systemd/user/restic-prune.timer
[Unit]
Description=prune restic snapshots

[Timer]
Persistent=true
OnCalendar=daily
Unit=restic-prune.service

[Install]
WantedBy=timers.target
