# restic/.config/systemd/user/restic-backup.service
[Unit]
Description=restic Backup
Documentation=man:restic(1)

[Service]
Type=simple
ExecStartPre=nm-online
ExecStart=%h/.dotfiles/bin/restic-backup
