# restic/.config/systemd/user/restic-backup.timer
[Unit]
Description=restic backup

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=timers.target
