# restic/.config/systemd/user/restic-backup.timer
[Unit]
Description=restic backup

[Timer]
Persistent=true
OnCalendar=hourly
Unit=restic-backup.service

[Install]
WantedBy=timers.target
