# restic/.config/systemd/user/restic-backup.timer
[Unit]
Description=restic Backup Trigger
Documentation=man:restic(1)

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=timers.target
