# restic/.config/systemd/user/restic-prune.timer
[Unit]
Description=restic prune snapshots

[Timer]
Persistent=true
OnCalendar=weekly
Unit=restic-prune.service

[Install]
WantedBy=timers.target